# wayland/.config/systemd/user/dunst.service.d/50-wayland.conf
[Unit]
PartOf=
PartOf=wayland-session.target
After=wayland-session.target

[Install]
WantedBy=wayland-session.target
