# r2e/.config/systemd/user/rss2email.service
[Unit]
Description=RSS email synchronization service

[Service]
Type=oneshot
ExecStart=/usr/bin/r2e run

[Install]
WantedBy=default.target
