# gruvbox/.config/profile.d/20-dircolors.sh
#!/bin/sh

eval $(dircolors "$XDG_CONFIG_HOME"/gruvbox.dircolors) 
