# .config/profile.d/20-gruvbox.sh
#!/bin/sh

export BASE16_THEME_DEFAULT="gruvbox-dark-hard"
export BASE16_SHELL_ENABLE_VARS=1
