# .config/profile.d/20-zk.sh
#!/bin/sh

export ZK_NOTEBOOK_DIR="$HOME/notebook"
