# .config/profile.d/50-gdb.sh
#!/bin/sh

export GDBHISTFILE="$XDG_DATA_HOME/gdb/history"
export DEBUGINFOD_URLS="https://debuginfod.elfutils.org/"
