# aerc/.config/systemd/user/maildir-rank-addr@.timer
[Unit]
Description=Addressbook generation timer

[Timer]
OnUnitActiveSec=12h
Unit=maildir-rank-addr@%i.service

[Install]
WantedBy=timers.target
