# i3status/.config/systemd/user/rss-status.service
[Unit]
Description=rss status bar block

[Service]
Type=oneshot
ExecStart=%h/.local/bin/rss-status
