# i3status/.config/systemd/user/rss-status.timer
[Unit]
Description=rss status bar block

[Timer]
OnActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
