# python/.config/profile.d/20-python.sh
#!/bin/sh

export PYTHONSTARTUP="$XDG_CONFIG_HOME/python/startup.py"
