# nvim/.config/profile.d/20-nvim.sh
#!/bin/sh

export EDITOR="nvim"
export DIFFPROG="$EDITOR -d"
