# ssh/.config/systemd/user/ssh-open@.service
[Unit]
Description=xdg-open socket listener

[Service]
ExecStart=/usr/bin/xargs swaymsg exec xdg-open
StandardInput=socket
StandardOutput=null
