# ssh/.config/systemd/user/ssh-open.socket
[Unit]
Description=xdg-open socket listener

[Socket]
ListenStream=%t/xdg-open.socket
Accept=yes

[Install]
WantedBy=sockets.target
