# sh/.config/profile.d/10-base16.sh
#!/bin/sh

BASE16_COLOR_00_HEX=${BASE16_COLOR_00_HEX:-1d2021}
BASE16_COLOR_01_HEX=${BASE16_COLOR_01_HEX:-3c3836}
BASE16_COLOR_02_HEX=${BASE16_COLOR_02_HEX:-504945}
BASE16_COLOR_03_HEX=${BASE16_COLOR_03_HEX:-665c54}
BASE16_COLOR_04_HEX=${BASE16_COLOR_04_HEX:-bdae93}
BASE16_COLOR_05_HEX=${BASE16_COLOR_05_HEX:-d5c4a1}
BASE16_COLOR_06_HEX=${BASE16_COLOR_06_HEX:-ebdbb2}
BASE16_COLOR_07_HEX=${BASE16_COLOR_07_HEX:-fbf1c7}
BASE16_COLOR_08_HEX=${BASE16_COLOR_08_HEX:-fb4934}
BASE16_COLOR_09_HEX=${BASE16_COLOR_09_HEX:-fe8019}
BASE16_COLOR_0A_HEX=${BASE16_COLOR_0A_HEX:-fabd2f}
BASE16_COLOR_0B_HEX=${BASE16_COLOR_0B_HEX:-b8bb26}
BASE16_COLOR_0C_HEX=${BASE16_COLOR_0C_HEX:-8ec07c}
BASE16_COLOR_0D_HEX=${BASE16_COLOR_0D_HEX:-83a598}
BASE16_COLOR_0E_HEX=${BASE16_COLOR_0E_HEX:-d3869b}
BASE16_COLOR_0F_HEX=${BASE16_COLOR_0F_HEX:-d65d0e}

export BASE16_COLOR_00_HEX BASE16_COLOR_01_HEX BASE16_COLOR_02_HEX BASE16_COLOR_03_HEX \
	BASE16_COLOR_04_HEX BASE16_COLOR_05_HEX BASE16_COLOR_06_HEX BASE16_COLOR_07_HEX \
	BASE16_COLOR_08_HEX BASE16_COLOR_09_HEX BASE16_COLOR_0A_HEX BASE16_COLOR_0B_HEX \
	BASE16_COLOR_0C_HEX BASE16_COLOR_0D_HEX BASE16_COLOR_0E_HEX BASE16_COLOR_0F_HEX
