# sh/.config/profile.d/10-systemd.sh
#!/bin/sh

export $(run-parts /usr/lib/systemd/user-environment-generators | xargs)
