# sh/.config/profile.d/20-distcc.sh
#!/bin/sh

export DISTCC_DIR="$XDG_CONFIG_HOME"/distcc
