# sh/.config/profile.d/20-dotnet.sh
#!/bin/sh

export DOTNET_CLI_HOME="$XDG_DATA_HOME/dotnet"
export PATH="$PATH:$DOTNET_CLI_HOME/tools"
