# sh/.config/profile.d/20-gruvbox.sh
#!/bin/sh

export BASE16_DEFAULT_THEME="gruvbox-dark-hard"
