# sh/.config/profile.d/20-pass.sh
#!/bin/sh

export PASSWORD_STORE_DIR="$HOME/.password-store"
export PASSWORD_STORE_ENABLE_EXTENSIONS=true
