# sh/.config/profile.d/20-wget.sh
#!/bin/sh

export WGETRC="$XDG_CONFIG_HOME/wgetrc"
