# sh/.config/profile.d/50-android.sh
#!/bin/sh

export ANDROID_USER_HOME="$XDG_DATA_HOME"/android
