# sh/.config/profile.d/50-opam.sh
#!/bin/sh

export OPAMROOT="$XDG_DATA_HOME/opam"
