# sh/.config/profile.d/50-ollama.sh
#!/bin/sh

export OLLAMA_MODELS="$XDG_DATA_HOME"/ollama/models
