# sh/.config/profile.d/50-platformio.sh
#!/bin/sh

export PLATFORMIO_CORE_DIR="$XDG_DATA_HOME/platformio"
