# sh/.config/profile.d/50-wine.sh
#!/bin/sh

export WINEPREFIX="$XDG_DATA_HOME"/wine
