# bemenu/.config/profile.d/50-bemenu.sh
#!/usr/bin/sh

base00="#1d2021"
base01="#3c3836"
base02="#504945"
base03="#665c54"
base04="#bdae93"
base05="#d5c4a1"
base06="#ebdbb2"
base07="#fbf1c7"
base08="#fb4934"
base09="#fe8019"
base0A="#fabd2f"
base0B="#b8bb26"
base0C="#8ec07c"
base0D="#83a598"
base0E="#d3869b"
base0F="#d65d0e"

export BEMENU_OPTS="--ignorecase \
--list 10 \
--scrollbar autohide \
--line-height 18 \
--binding vim \
--fn 'Fira Code Nerd Font 12' \
--no-overlap \
--ab '$base02' \
--af '$base06' \
--fb '$base00' \
--ff '$base06' \
--scb '$base06' \
--scf '$base00' \
--nb '$base02' \
--nf '$base06' \
--tb '$base0D' \
--tf '$base06' \
--hb '$base03' \
--hf '$base06' \
--sb '$base06' \
--sf '$base06'"
