# notmuch/.config/profile.d/50-notmuch.sh
#!/bin/sh

export NOTMUCH_CONFIG="$XDG_CONFIG_HOME"/notmuch/default/config
